# src/aiguard/installer/templates/ai-guard.socket.in
[Unit]
Description=Datadog AI Guard proxy listening socket

[Socket]
ListenStream=${LISTEN_STREAM}
Service=${SERVICE_NAME}

[Install]
WantedBy=sockets.target
