# src/aiguard/installer/templates/ai-guard.service.in
[Unit]
Description=Datadog AI Guard proxy for coding agents
After=network.target

[Service]
Type=simple
ExecStart=${WRAPPER}
Restart=on-failure
RestartSec=2
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ai-guard

[Install]
WantedBy=default.target
